# deploy/macmini/com.powerplanner.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- 일일그래프: 매일 06:00 (KST) 어제 하루 24시간 사용량 그래프 전송.
     __PROJECT_DIR__ 3곳을 실제 클론 경로로 치환한 뒤 ~/Library/LaunchAgents/ 에 설치. -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.powerplanner.daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>__PROJECT_DIR__/deploy/macmini/run.sh</string>
    <string>daily_report.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__PROJECT_DIR__</string>

  <key>StartCalendarInterval</key>
  <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>

  <key>StandardOutPath</key>
  <string>__PROJECT_DIR__/logs/launchd.daily.log</string>
  <key>StandardErrorPath</key>
  <string>__PROJECT_DIR__/logs/launchd.daily.log</string>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
